# root/etc/services.d/hydroxide/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc hydroxide $HYDROXIDE_CMD $HYDROXIDE_FLAGS
